# pingstate.service
[Unit]
Description=pingstate — TCP probe state machine with local HTTP API
Documentation=https://github.com/pacnpal/pingstate
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
# Adjust the protocol/address/port/check and the ExecStart path for your setup.
# After `pip install pingstate` the console script lives at /usr/local/bin/pingstate
# (or wherever your Python install puts scripts).
ExecStart=/usr/local/bin/pingstate \
    --protocol tcp \
    --address 1.1.1.1 \
    --port 443 \
    --interval 5 \
    --api-host 127.0.0.1 \
    --api-port 8787
Restart=on-failure
RestartSec=5

# Run as an unprivileged user. Create one first, e.g.:
#   sudo useradd --system --no-create-home --shell /usr/sbin/nologin pingstate
User=pingstate
Group=pingstate

# Hardening. Safe defaults for a stdlib-only daemon that only opens sockets.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true

[Install]
WantedBy=multi-user.target
